# deployment/distconv-transcoding-engine.service
[Unit]
Description=DistConv Transcoding Engine - Video Processing Worker
Documentation=https://github.com/segin/distconv
After=network.target network-online.target distconv-dispatch.service
Wants=network-online.target
RequiresMountsFor=/var/lib/distconv /opt/distconv

[Service]
Type=simple
User=distconv
Group=distconv
WorkingDirectory=/var/lib/distconv/engine

# Use FHS-compliant configuration
ExecStart=/opt/distconv/engine/bin/transcoding_engine_modern \
  --config /etc/distconv/transcoding-engine.conf

# Service lifecycle management
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=60

# Restart policy
Restart=always
RestartSec=15
StartLimitInterval=120
StartLimitBurst=5

# Health monitoring
WatchdogSec=60

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=no  # FFmpeg needs executable memory

# Writable directories (FHS-compliant)
ReadWritePaths=/var/lib/distconv/engine
ReadWritePaths=/var/log/distconv
ReadWritePaths=/run/distconv
ReadWritePaths=/tmp/distconv
ReadWritePaths=/var/cache/distconv

# Capability restrictions (FFmpeg needs some capabilities)
CapabilityBoundingSet=CAP_SETUID CAP_SETGID
AmbientCapabilities=

# System call filtering (more permissive for FFmpeg)
SystemCallFilter=@system-service
SystemCallFilter=~@debug
SystemCallFilter=~@mount
SystemCallFilter=~@reboot
SystemCallFilter=~@swap
SystemCallErrorNumber=EPERM

# Network restrictions
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Resource limits (higher for transcoding)
LimitNOFILE=65536
LimitNPROC=512
MemoryLimit=8G
CPUQuota=400%

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=distconv-transcoding-engine
SyslogFacility=daemon
SyslogLevel=info

# Performance optimizations for media processing
Nice=-10
IOSchedulingClass=1
IOSchedulingPriority=2
CPUSchedulingPolicy=1

# Environment variables
Environment=TMPDIR=/tmp/distconv
Environment=FFMPEG_LOG_LEVEL=warning

[Install]
WantedBy=multi-user.target
Alias=distconv-engine.service